# rc/rc.local
source ~/.vim/dotvim/rc/vimrc

try
    if (has('gui_running'))
        call SwitchTheme(0)
    else
        "colo wombat256
        highlight CursorLine cterm=NONE ctermbg=NONE
    endif
    "call RandomThemes(1)
catch
endtry


if(has('gui_macvim'))
    "set transparency = 40
endif

