# rc/rc.local
source ~/.vim/dotvim/rc/vimrc

try
    if (has('gui_running'))
        call SwitchTheme(0)
        highlight clojureError term=undercurl ctermfg=15 ctermbg=4 gui=undercurl guifg=White guibg=Red
    else
        "colo wombat256
        highlight CursorLine cterm=NONE ctermbg=NONE
    endif
    "call RandomThemes(1)
catch
endtry


if(has('gui_macvim'))
    "set transparency=20
endif

