# rc/rc.local
source ~/.vim/dotvim/rc/vimrc

try
    if (has('gui_running'))
        if(has('gui_macvim'))
            call SwitchTheme(26)
        else
            "call SwitchTheme(26)
            call SwitchTheme(7)
        endif
    else
        "colo wombat256
        highlight CursorLine cterm=NONE ctermbg=NONE
    endif
    "call RandomThemes(1)
catch
endtry


if(has('gui_macvim'))
    set transparency=0
endif
