# rc/rc.local
source ~/.vim/dotvim/rc/vimrc

try
    if (has('gui_running'))
        call SwitchTheme(7)
    else
        call SwitchTheme(0)
    endif
    "call RandomThemes(1)
catch
endtry


if(has('gui_macvim'))
    "set transparency=20
endif

